# embedded-system-deployer.service
[Unit]
Description=Embedded System Deployer (from embedded-system-manager: https://github.com/b-mer/embedded-system-manager)
After=network.target systemd-user-sessions.service
Conflicts=getty@tty1.service

[Service]
User=root
Group=root
Type=idle
RuntimeDirectory=embedded-system-manager
RuntimeDirectoryMode=0700
ExecStart=/opt/embedded-system-manager/startup.sh
Restart=on-failure
RestartSec=10
IgnoreSIGHUP=yes
StandardInput=tty
StandardOutput=journal
StandardError=journal
TTYPath=/dev/tty1
TTYReset=yes
TTYVHangup=yes
TTYVTDisallocate=yes

# Environment variables for Wayland/Cage
Environment="XDG_RUNTIME_DIR=/run/embedded-system-manager"
Environment="WLR_LIBINPUT_NO_DEVICES=1"

[Install]
WantedBy=multi-user.target
